# deploy/systemd/systemmanager-mcp.service
[Unit]
Description=SystemManager MCP Server with OAuth/OIDC
Documentation=https://github.com/mdlmarkham/SystemManager
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/systemmanager

# Load secrets from protected environment file
EnvironmentFile=/opt/systemmanager/.env

# Use the virtual environment Python
ExecStart=/opt/systemmanager/venv/bin/python -m src.mcp_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
StartLimitInterval=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/systemmanager/logs

# Resource limits
LimitNOFILE=65536

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=systemmanager-mcp

[Install]
WantedBy=multi-user.target